# .config/systemd/user/do-backup.service
[Unit]
Description=Create a backup for the month

[Service]
Type=oneshot
ExecStart=/bin/bash /home/mssola/bin/do-backup.sh